# Section_1/crontab
0 * * * * /usr/bin/python3 /path/to/section_1.py